# crontab
# ============================================================
# Scraper Schedule — Server 187.77.168.42
# All times UTC — staggered to avoid overlap
# ============================================================

# Bien'ici (France) — 03:00
0 3 * * *   /app/scheduler.sh bienici >> /app/logs/cron.log 2>&1

# Mubawab (Tunisia) — 04:00
0 4 * * *   /app/scheduler.sh mubawab >> /app/logs/cron.log 2>&1

# MktList (Canada) — 05:00
0 5 * * *   /app/scheduler.sh mktlist >> /app/logs/cron.log 2>&1

# PropertyFinder (Egypt) — 06:00
0 6 * * *   /app/scheduler.sh propertyfinder >> /app/logs/cron.log 2>&1

# Sync MongoDB -> PostgreSQL — 10:00 (after scrapers finish)
0 10 * * *  /app/scheduler.sh sync >> /app/logs/cron.log 2>&1

# Health monitor — every 6h
0 */6 * * * /app/scheduler.sh monitor >> /app/logs/cron.log 2>&1